# templates/freshclam.timer
[Unit]
Description=update clamav definition
Wants=freshclam.service

[Timer]
OnCalendar=*-*-* 08:00:00

[Install]
WantedBy=timers.target
